# opuntia_imagestream/package/compex-production/files/etc/uci-defaults/50-prod_info
#!/bin/bash
#
# Copyright (C) 2013-2014 ImageStream Internet Solutions, Inc.
#
# One-shot production info text file for Compex Production line
# File will disappear after reboot

cp -a /root/prod_info.sh /tmp
/bin/bash /tmp/prod_info.sh &
