# opuntia_imagestream/package/envoy-configs/files/etc/uci-defaults/12_dhcp
#!/bin/sh
#
# Copyright (C) 2013 ImageStream Internet Solutions, Inc.
#

exit 0


