# overlay/base-files/files/etc/hotplug.d/net/40-default-qdisc-cake
#!/bin/sh

if [ "$ACTION" = add ]; then
	tc qdisc replace dev $DEVICENAME root cake
fi
