# overlay/base-files/files/etc/hotplug.d/net/40-wifi-default-qdisc-cake
#!/bin/sh

if [ "$ACTION" == add ] && echo "$DEVICENAME" | grep '^wlan'; then
	tc qdisc replace dev $DEVICENAME root fq_codel
fi
