# opuntia_imagestream/package/eudev/files/udev.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=10
STOP=90

boot() {
        /sbin/udevadm hwdb -u
        /sbin/udevd -d
}
