# overlay/quagga/files/etc/uci-defaults/15_quagga-fixconf
#!/bin/sh

mv /etc/quagga /etc/quagga.bak
cp -a /etc/quagga.bak /etc/quagga
rm -rf /etc/quagga.bak
