# target/linux/x86/base-files/etc/init.d/isis-postinit
#!/bin/sh /etc/rc.common
#
#

START=99

start() {
	local rootfsdev
	local rootfstype

	. /lib/functions.sh
	. /lib/upgrade/platform.sh

	RC_POSTINSTALL="rc.postinstall"
	mnt_nand_root=/tmp/nand_root
	mkdir $mnt_nand_root
	
	rootfstype="$(awk 'BEGIN { RS=" "; FS="="; } ($1 == "rootfstype") { print $2 }' < /proc/cmdline)"
	case "$rootfstype" in
		squashfs|jffs2)
			rootfsdev="$(awk 'BEGIN { RS=" "; FS="="; } ($1 == "block2mtd.block2mtd") { print substr($2,1,index($2, ",")-1) }' < /proc/cmdline)";;
		ext4)
			rootfsdev="$(awk 'BEGIN { RS=" "; FS="="; } ($1 == "root") { print $2 }' < /proc/cmdline)";;
	esac

	mount -t ext4 -o rw,noatime "${rootfsdev%[0-9]}1" $mnt_nand_root

	if [ -x ${mnt_nand_root}/${RC_POSTINSTALL} ] ; then
		logger "Running postinstall script..."
		${mnt_nand_root}/${RC_POSTINSTALL}
		logger "Done running postinstall script."
	fi
	umount $mnt_nand_root
	rmdir $mnt_nand_root
}
