# target/linux/x86/base-files/etc/init.d/watchdog
#!/bin/sh /etc/rc.common
#
#

START=97

start() {
	/sbin/watchdog /dev/watchdog
}

stop() {
	killall watchdog
}
