# overlay/base-files/files/etc/uci-defaults/00_haveged
#!/bin/sh

/etc/init.d/haveged start
exit 0
