# imagestream/package/netifd-qmi/files.old/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

chmod u+x /etc/savevar
/etc/savevar 

uci set wireless.@wifi-device[0].disabled=0; uci commit wireless; wifi

exit 0
